# deploy/systemd/loan-assistant-celery-heavy.service
[Unit]
Description=Loan Assistant Celery Heavy Worker
After=network.target redis-server.service
Requires=redis-server.service

[Service]
User=root
WorkingDirectory=/root/desktop_app
EnvironmentFile=/root/desktop_app/.env
ExecStart=/root/desktop_app/venv/bin/celery -A backend.celery_app.celery_app worker -l info -Q heavy --concurrency=1 --prefetch-multiplier=1 --max-tasks-per-child=10 -n heavy@%h
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
